# systemd/user/battery_notify.service
[Unit]
Description=Battery notification service

[Service]
Type=oneshot
ExecStart=/home/heisenberg/.config/hypr/battery-info.sh
